# files/etc/supervisord.d/app_tasks.ini
[program:app_tasks]
user=commit--blog
directory=/var/commit--blog
command=/var/commit--blog/venv/bin/python manage.py run_tasks
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/app_tasks.log
stdout_logfile_maxbytes = 40MB
stdout_logfile_backups = 2
