# files/etc/supervisord.d/gunicorn.ini
[program:gunicorn]
user=commit--blog
directory=/var/commit--blog
command=/var/commit--blog/venv/bin/newrelic-admin run-program venv/bin/gunicorn --no-sendfile --workers 2 --worker-class gevent --keep-alive 5 --bind 127.0.0.1:9000 wsgi:app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/gunicorn.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10
